# debian/DEBIAN/postinst
#!/bin/bash

echo "Start install hobot audio config"
systemctl daemon-reload
systemctl enable hobot-audio.service
exit 0

