# debian/DEBIAN/postrm
#!/bin/bash

sudo systemctl disable example_app.service
exit 0