# debian/DEBIAN/postrm
#!/bin/bash

sudo systemctl disable hobot-audio.service
exit 0